# nodepilot-frontend.service
[Unit]
Description=NodePilot Frontend
Documentation=https://github.com/yourrepo/nodepilot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/deployer/frontend
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nodepilot-frontend

[Install]
WantedBy=multi-user.target
